# etc/local.d/killplymouthd.start
#!/bin/sh
# Show below message on system startup
pkill plymouthd
